# scripts/resources/agents/agent-s2/docker/config/supervisor.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
user=agents2
autorestart=true
priority=100
stdout_logfile=/tmp/xvfb.log
stderr_logfile=/tmp/xvfb.err
environment=HOME="/home/agents2",USER="agents2"

[program:fluxbox]
command=/usr/bin/fluxbox -display :99
user=agents2
environment=DISPLAY=":99",HOME="/home/agents2",USER="agents2"
autorestart=true
priority=200
stdout_logfile=/tmp/fluxbox.log
stderr_logfile=/tmp/fluxbox.err
startsecs=5

[program:vnc-server]
command=/usr/bin/x11vnc -forever -usepw -shared -rfbport 5900 -display :99 -noxrecord -noxfixes -noxdamage -nomodtweak -repeat
user=agents2
environment=DISPLAY=":99",HOME="/home/agents2",USER="agents2"
autorestart=true
priority=300
stdout_logfile=/tmp/vnc.log
stderr_logfile=/tmp/vnc.err
startsecs=5

[program:startup]
command=/opt/agent-s2/startup.sh
user=agents2
autorestart=false
priority=350
stdout_logfile=/tmp/startup.log
stderr_logfile=/tmp/startup.err
startsecs=0

[program:security-proxy]
command=/home/agents2/.local/bin/mitmdump -s /opt/agent-s2/agent_s2/server/services/proxy_addon.py --mode transparent --listen-port 8085 --set confdir=/home/agents2/.mitmproxy --set block_global=false --set upstream_cert=false --set connection_strategy=lazy --set stream_large_bodies=50m --set body_size_limit=50m --set keep_host_header=true
user=agents2
directory=/opt/agent-s2
environment=DISPLAY=":99",PYTHONPATH="/opt/agent-s2:/home/agents2/.local/lib/python3.11/site-packages",HOME="/home/agents2",USER="agents2",AGENT_S2_SECURITY_PROFILE="moderate"
autorestart=true
priority=375
stdout_logfile=/tmp/proxy.log
stderr_logfile=/tmp/proxy.err
startsecs=10
stopwaitsecs=15

[program:firefox-monitor]
command=/opt/agent-s2/firefox-monitor.sh
user=agents2
directory=/opt/agent-s2
environment=DISPLAY=":99",HOME="/home/agents2",USER="agents2",MOZ_CRASHREPORTER_DISABLE="1",MOZ_DISABLE_CONTENT_SANDBOX="1"
autorestart=true
priority=380
stdout_logfile=/tmp/firefox-monitor.log
stderr_logfile=/tmp/firefox-monitor.err
startsecs=10
stopwaitsecs=15

[program:agent-s2-api]
command=python -m agent_s2.server
user=agents2
directory=/opt/agent-s2
environment=DISPLAY=":99",PYTHONPATH="/home/agents2/.local/lib/python3.11/site-packages",HOME="/home/agents2",USER="agents2",XAUTHORITY="/home/agents2/.Xauthority"
autorestart=true
priority=400
stdout_logfile=/tmp/api.log
stderr_logfile=/tmp/api.err
startsecs=10
stopwaitsecs=30

[group:agent-s2]
programs=xvfb,fluxbox,vnc-server,startup,security-proxy,firefox-monitor,agent-s2-api
priority=999